# kde4-profiles/netbook/share/autostart/plasma-netbook.desktop
[Desktop Entry]
Exec=plasma-netbook
Hidden=false
OnlyShowIn=KDE;
X-KDE-autostart-phase=BaseDesktop
